# crontab
3 * * * * run-parts /home/frew/.crontab.d/hourly
