# crontab
3 * * * * . "$HOME/.env" && run-parts "$HOME/.crontab.d/hourly"
* * * * * . "$HOME/.env" && run-parts "$HOME/.crontab.d/minutely"
